# redis.conf
# Upstart file at /etc/init/redis.conf
# Redis

description "Redis"

start on runlevel [23]
stop on shutdown
expect daemon
respawn

exec sudo -u redis /usr/bin/redis-server /etc/redis/redis.conf
